# systemd/claude-agent.service
[Unit]
Description=Claude Code Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/agent
EnvironmentFile=/home/ubuntu/agent/.env
Environment=TZ=Australia/Melbourne
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5
OOMPolicy=stop
MemoryMax=1536M
MemoryHigh=1280M
TasksMax=50
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/workspace /home/ubuntu/.claude-agent /home/ubuntu/agent /home/ubuntu/.claude /home/ubuntu/.config
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
